# deploy/systemd/esports-site-selection.service
[Unit]
Description=Esports Site Selection FastAPI Service
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=esports-site-selection
Group=esports-site-selection
WorkingDirectory=__APP_ROOT__/backend
EnvironmentFile=/etc/esports-site-selection/backend.env
Environment=PYTHONUNBUFFERED=1
ExecStart=__APP_ROOT__/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/esports-site-selection

[Install]
WantedBy=multi-user.target
